# systemd-userlevel/simple.timer
[Unit]
Description=Timer for periodically running simple bash script at user level

[Timer]
OnUnitActiveSec=60
