# systemd-systemlevel/github.timer
[Unit]
Description=Timer for periodically removing github followers

[Timer]
OnUnitActiveSec=60
